# debian/postinst
#!/bin/bash
if [ "$(cat /etc/modules | grep sensors-detect)" = "" ]; then
    /usr/bin/yes| sudo sensors-detect; sudo /etc/init.d/lm-sensors stop;sudo update-rc.d -f lm-sensors remove	
fi
